# gunicorn.conf.j2
[program:gunicorn]
command=/usr/local/bin/gunicorn pgw:app -c /home/{{ username }}/gunicorn.conf.py
directory=/home/{{ username }}
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
